# root/etc/s6-overlay/s6-rc.d/svc-trproxy/run
#!/usr/bin/with-contenv bash

cd /config

exec s6-setuidgid abc /usr/sbin/trproxy -conf /config
