# root/etc/s6-overlay/s6-rc.d/init-trproxy-config/run
#!/usr/bin/with-contenv bash

mkdir -p /config/trproxy

# permissions
chown abc:abc -R /config/trproxy

setcap cap_net_admin,cap_net_raw=eip /usr/sbin/trproxy
